# systemd/security-manager-cleanup.service.in
[Unit]
Description=Cleanup private sharing
DefaultDependencies=no
After=sysinit.target
Requires=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart = @BIN_INSTALL_DIR@/security-manager-cleanup
SmackProcessLabel=System::Privileged

