# systemd/security-server.target
[Unit]
Description=security-server sockets
DefaultDependencies=true

