# systemd/security-manager-slave.service.in
[Unit]
Description=Start the security manager slave
ConditionVirtualization=lxc

[Service]
Type=notify
ExecStart=@BIN_INSTALL_DIR@/security-manager --slave

Sockets=security-manager-slave.socket
