# systemd/security-manager.service.in
[Unit]
Description=Start the security manager
Before=dbus.service

[Service]
Type=notify
ExecStart=@BIN_INSTALL_DIR@/security-manager
Sockets=security-manager.socket
Restart=always
SmackProcessLabel=System::Privileged

[Install]
WantedBy=dbus.service
