# systemd/security-manager.target
[Unit]
Description=security-manager sockets
DefaultDependencies=true

