# systemd/security-manager.service.in
[Unit]
Description=Start the security manager

[Service]
Type=notify
ExecStart=@BIN_INSTALL_DIR@/security-manager
Sockets=security-manager.socket
