# systemd/security-server-password-check.socket
[Socket]
ListenStream=/tmp/.security-server-api-password-check.sock
SocketMode=0777
SmackLabelIPIn=security-server::api-password-check
SmackLabelIPOut=@

Service=security-server.service

[Install]
WantedBy=sockets.target
