# packaging/security-server.service

[Unit]
Description=Start the security server

[Service]
ExecStart=/usr/bin/security-server

[Install]
WantedBy=multi-user.target
